# contrib/init/californium252d.conf
description "Californium Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CALIFORNIUMD_BIN="/usr/bin/californium252d"
env CALIFORNIUMD_USER="californium252"
env CALIFORNIUMD_GROUP="californium252"
env CALIFORNIUMD_PIDDIR="/var/run/californium252d"
# upstart can't handle variables constructed with other variables
env CALIFORNIUMD_PIDFILE="/var/run/californium252d/californium252d.pid"
env CALIFORNIUMD_CONFIGFILE="/etc/californium252/californium252.conf"
env CALIFORNIUMD_DATADIR="/var/lib/californium252d"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # californium252d will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CALIFORNIUMD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run californium252d."
        echo "The setting must appear in $CALIFORNIUMD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Californium Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CALIFORNIUMD_PIDDIR"
    chmod 0755 "$CALIFORNIUMD_PIDDIR"
    chown $CALIFORNIUMD_USER:$CALIFORNIUMD_GROUP "$CALIFORNIUMD_PIDDIR"
    chown $CALIFORNIUMD_USER:$CALIFORNIUMD_GROUP "$CALIFORNIUMD_CONFIGFILE"
    chmod 0660 "$CALIFORNIUMD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CALIFORNIUMD_PIDFILE" \
    --chuid $CALIFORNIUMD_USER:$CALIFORNIUMD_GROUP \
    --exec "$CALIFORNIUMD_BIN" \
    -- \
    -pid="$CALIFORNIUMD_PIDFILE" \
    -conf="$CALIFORNIUMD_CONFIGFILE" \
    -datadir="$CALIFORNIUMD_DATADIR" \
    -disablewallet \
    -daemon

